# me909s/files/etc/init.d/me909s
#!/bin/sh /etc/rc.common
# Start Huawei modem

START=99
STOP=01

USE_PROCD=1

start_service() {
    procd_open_instance
    procd_set_param command sh "usr/bin/me909s.sh"
    procd_close_instance
}
